# vendor/shieldforce/scoob/progs/runs/laravel/8.2/supervisord/conf/mariadb.conf
### to schedule run, uncomment all
[program:mariadb]
command=/usr/sbin/mysqld
nodaemon=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_maxbytes=0
stderr_maxbytes=0
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
autostart=true
autorestart=false
startretries = 1
user=root



