# vendor/shieldforce/scoob/progs/runs/laravel/8.3/supervisord/conf/redis.conf
### to schedule run, uncomment all
[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
autorestart=false
user=root
stdout_logfile=/var/log/redis/stdout.log
stderr_logfile=/var/log/redis/stderr.log



