# vendor/shieldforce/scoob/progs/runs/laravel/7.3/supervisord/conf/horizon.conf
### to schedule run, uncomment all
[program:horizon]
process_name=%(program_name)s
command=php /var/www/artisan horizon
autostart=true
autorestart=false
user=root
redirect_stderr=true
stdout_logfile=/var/www/horizon.log
stopwaitsecs=3600



